# skwhisper.service
[Unit]
Description=SKWhisper — Sovereign Subconscious Memory Layer (%I)
After=network.target

[Service]
Type=simple
# Instance name = agent name (e.g. systemctl --user start skwhisper@lumina)
Environment=SKCAPSTONE_AGENT=%i
ExecStart=/usr/bin/python3 -m skwhisper daemon
WorkingDirectory=%h/clawd/projects/skwhisper
Environment=PYTHONPATH=%h/clawd/projects/skwhisper
Environment=PATH=%h/.skenv/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=30

# Logging — routed to agent-specific dir via SKCAPSTONE_AGENT
StandardOutput=append:%h/.skcapstone/agents/%i/skwhisper/daemon.log
StandardError=append:%h/.skcapstone/agents/%i/skwhisper/daemon.log

[Install]
WantedBy=default.target
